# projects/scripts/baimeng-cron.timer
[Unit]
Description=白梦写作网定时任务 - 每天9点执行

[Timer]
OnCalendar=*-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
